# templates/inventory-api@.service.epp
<%- | | -%>
# Managed by Puppet
[Unit]
Description=Inventory API Server (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inventory-%i
Group=inventory-%i
EnvironmentFile=/etc/inventory-system/%i.conf
Environment="ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY}"
ExecStart=/usr/local/bin/inventory-md api ${INVENTORY_PATH} --port ${API_PORT} --host ${API_HOST}
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventory-api-%i

[Install]
WantedBy=multi-user.target
